# mass-compress.conf
description "start mass-compress server"

start on ( desktop-session-start )
stop on runlevel[016]
respawn

exec /usr/local/bin/mass-compress
